# l33t-0.1/debian/l33t.cron.d
37 13 * * * root    /bin/bash /opt/1337.sh -e
38 13 * * * root    /bin/bash /opt/1337.sh -d
